# systemd/ai-project-backup.service
[Unit]
Description=بکاپ روزانه پروژه هوش مصنوعی

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ai-project-backup.sh
