# systemd/ai-project-backup.timer
[Unit]
Description=اجرای روزانه بکاپ پروژه هوش مصنوعی

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target